# newsfeed-ai-bot.service
[Unit]
Description=Newsfeed AI Discord Bot
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/newsfeed-ai
ExecStart=/home/exedev/.bun/bin/bun run src/bot.ts
Restart=always
RestartSec=10
EnvironmentFile=/home/exedev/newsfeed-ai/.env

[Install]
WantedBy=multi-user.target
