# systemd/newsfeed-ai-bot.service
[Unit]
Description=Newsfeed AI Discord Bot
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=/home/__USER__/newsfeed-ai
ExecStart=/home/__USER__/.bun/bin/bun run apps/bot/src/main.ts
Restart=always
RestartSec=10
EnvironmentFile=/home/__USER__/newsfeed-ai/.env

[Install]
WantedBy=multi-user.target
